# example/crontab
0 2 * * 0 /scripts/backup.sh >> /data/backup.log 2>&1
